# 60-hdmi-switch-config.rules
ACTION=="add|change",SUBSYSTEM=="extcon",ENV{DEVPATH}=="/devices/platform/soc/6000000.hdmi/extcon/hdmi",RUN+="/usr/bin/batocera-hdmi-switch-config"
